# deploy/systemd/codebot-frontend.service
[Unit]
Description=CodeBot Next.js Frontend
After=network.target aicodebot.service

[Service]
Type=simple
User=omatic657
WorkingDirectory=/home/omatic657/aicoderbot/apps/codebot-builder
Environment="PORT=3000"
Environment="NODE_ENV=production"
ExecStart=/usr/bin/pnpm start
Restart=on-failure
RestartSec=5
StandardOutput=append:/tmp/nextjs.log
StandardError=append:/tmp/nextjs.log

[Install]
WantedBy=multi-user.target
